# crontab
30 6 * * * export TZ=Europe/Berlin; /app/blitzerbot -prodmode true > /proc/1/fd/1 2>/proc/1/fd/2
