# crontab
SHELL=/bin/bash
BASH_ENV=/etc/environment
30 6 * * * /app/blitzerbot -prodmode true >/proc/1/fd/1 2>/proc/1/fd/2
